# .config/systemd/user/sync_work_calendar.service
[Unit]
Description=Sync work calendar to org agenda

[Service]
ExecStart=/home/omad/bin/sync-work-calendar
